# templates/tmux/tmux_autosave.timer
[Unit]
Description=tmux autosave every 5 min

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
